# .config/autostart/picom.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/picom --experimental-backends
NoDisplay=true
OnlyShowIn=bspwm
Icon=picom
Name=picom
