# .config/autostart/pipewire.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/pipewire
NoDisplay=true
OnlyShowIn=bspwm
Icon=pipewire
Name=pipewire
