# .config/autostart/xplugd.desktop
[Desktop Entry]
Type=Application
Exec=xplugd
NoDisplay=true
OnlyShowIn=bspwm
Name=xplugd
