# autostart/dunst.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/dunst
NoDisplay=true
OnlyShowIn=bspwm
Icon=dunst
Name=dunst
