# autostart/numlock.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/numlockx on
NoDisplay=true
OnlyShowIn=bspwm
Icon=numlockx
Name=numlockx
