# autostart/set_cursor.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/xsetroot -cursor_name left_ptr
NoDisplay=true
OnlyShowIn=bspwm
Icon=xsetroot
Name=Set cursor
