# autostart/sxhkd.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/sh -c "exec /usr/bin/sxhkd -c $HOME/.config/bspwm/sxhkdrc"
NoDisplay=true
OnlyShowIn=bspwm
Icon=sxhkd
Name=sxhkd
